# dist/lib/systemd/system/sphinkter.service
[Unit]
Description=Sphinkter
After=network.target

[Service]
ExecStart=/usr/share/sphinkter/sphincterd.py
Restart=on-failure

[Install]
Alias=sphinkter.service

